# DEBIAN/postinst
#!/bin/bash


. /usr/share/debconf/confmodule

package_name="restore-ee-frontend"
configdir="/etc/restore"
includefile="${configdir}/apache.conf"

case "$1" in
    configure)
	webserver="apache2"
	webserver=`echo $webserver|sed -e 's/,  */ /g'`
	for i in $webserver; do
	   	if [ ! -d /etc/$i/conf.d/ ]; then
			install -d -m755 /etc/$i/conf.d/
	    	fi
	    	if [ ! -e /etc/$i/conf.d/restore.conf ]; then
			ln -s ${includefile} /etc/$i/conf.d/restore.conf
			restart="$i $restart"
	   	fi
	done
	addgroup --system restore 2>/dev/null
	adduser www-data restore

	chown -R restore /var/lib/restore 2>/dev/null ||
	adduser --system --home /var/lib/restore --no-create-home --disabled-password --ingroup restore restore
	
	touch /var/log/restore/fastcgi.crash.log
	chown www-data:www-data /var/log/restore/fastcgi.crash.log
	
      	a2enmod fcgid
	a2enmod rewrite
	servers="apache2"
	. /usr/share/wwwconfig-common/restart.sh
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
         echo "postinst called with unknown argument \`$1'" >&2
         exit 1
    ;;
esac
db_stop || true

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installinit
if [ -x "/etc/init.d/restore-frontend" ]; then
	update-rc.d restore-frontend defaults >/dev/null
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d restore-frontend start || exit $?
	else
		/etc/init.d/restore-frontend start || exit $?
	fi
fi
# End automatically added section


exit 0
